# dart-tournament-web.service
[Unit]
Description=Dart Tournament Web App
After=network.target

[Service]
Type=simple
# Change both paths if you deploy elsewhere (e.g. /home/user/DartTournamentWebApp).
WorkingDirectory=/opt/DartTournamentWebApp
ExecStart=/opt/DartTournamentWebApp/target/release/web

# Listen on all interfaces, port 8080 (override with HOST/PORT if needed)
Environment=HOST=0.0.0.0
Environment=PORT=8080

# If you get status=203/EXEC: binary missing or not executable. Fix:
#   1. Build: cargo build --release --bin web
#   2. chmod +x /opt/DartTournamentWebApp/target/release/web

Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dart-tournament-web

[Install]
WantedBy=multi-user.target
